# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:bauplan-backend]
command=python3 main.py
directory=/app/backend
autostart=true
autorestart=true
startsecs=10
startretries=3
redirect_stderr=true
stdout_logfile=/app/logs/backend.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
environment=ENVIRONMENT=production,PYTHONPATH=/app/backend

[program:bauplan-frontend]
command=npm start
directory=/app/frontend
autostart=true
autorestart=true
startsecs=10
startretries=3
redirect_stderr=true
stdout_logfile=/app/logs/frontend.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
environment=NODE_ENV=production 